# deploy/systemd/opensoul-nspawn@.service
[Unit]
Description=OpenSoul Gateway nspawn instance %i
After=network.target

[Service]
ExecStart=/usr/bin/systemd-nspawn -D /var/lib/opensoul/%i/rootfs --private-users=pick --private-network --bind=/var/lib/opensoul/%i/data:/data --setenv=OPENSOUL_STATE_DIR=/data
KillMode=mixed
Restart=on-failure

[Install]
WantedBy=multi-user.target
